# data/amd-bc-250-gpu-governor.service.in
[Unit]
Description=Oberon GPU governor

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/amd-bc-250-gpu-governor
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
